# systemd/citedy-reddit-writer.timer
[Unit]
Description=Run Citedy Reddit writer 3× daily

[Timer]
OnCalendar=*-*-* 08:00:00
OnCalendar=*-*-* 14:00:00
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
